# rc.d/si4.in
#!/bin/sh
#
# $Id$
#
# PROVIDE: si4
# REQUIRE: LOGIN

. /etc/rc.subr

name="si4"
rcvar="si4_enable"

pidfile="@APP_RUNDIR@/app.pid"
command="@prefix@/sbin/${name}"
command_args=""
procname="ruby /usr/local/sbin/${name}"

load_rc_config ${name}

: ${si4_enable:="NO"}

run_rc_command "$1"
#EOF
