# si4.service.in
[Unit]
Description=si4

[Service]
Type=forking
PIDFile=@APP_RUNDIR@/si4.pid
ExecStart=@prefix@/sbin/si4
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target

